# deploy/systemd/hiddenoasis-accounting-frontend.service
[Unit]
Description=Hidden Oasis Accounting Frontend
After=network.target hiddenoasis-accounting-backend.service

[Service]
Type=simple
User=hiddenoasis
Group=hiddenoasis
WorkingDirectory=/opt/accounting-program-online/frontend
Environment=NODE_ENV=production
EnvironmentFile=/etc/hiddenoasis/accounting-frontend.env
ExecStart=/usr/bin/npm run start -- -H 127.0.0.1 -p 3000
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
